# 86-dmc.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/sbin/dmc-load"
